# templates/silph.crontab.j2
# Managed by Juju - edit at your own risk

30 * * * * www-data /srv/web/app/Console/cake -app /srv/web/app generateCatalogCache
